# example/rahunas.init.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:          rahunas
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Should-Start:      $local_fs $network
# Should-Stop:       $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Handles the rahunas daemon
# Description:       Starts, stops the rahunas daemon. 
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin:@prefix@/sbin:@prefix@/bin
prefix=@prefix@
exec_prefix=@exec_prefix@

DAEMON=@bindir@/rahunasd
INIT=@sysconfdir@/default/rahunas
FIREWALL=@sysconfdir@/rahunas/firewall.sh
WEBLOGIN_CONFIG=@sysconfdir@/rahunas/weblogin-config-update.sh
NAME=rahunasd
DESC="RahuNAS: Rahu Network Access Server"


# This value is overriden in /etc/default/rahunas if it exists
RUN_DAEMON=no

test -f $DAEMON || exit 0
test -f $FIREWALL || exit 0

set -e

test -f $INIT || exit 0 
. $INIT 

case "$1" in
  start)
	echo -n "Starting $DESC: "

  $FIREWALL start
  $WEBLOGIN_CONFIG

	if [ "$RUN_DAEMON" = "yes" ]; then
		start-stop-daemon --start --quiet --pidfile @localstatedir@/run/$NAME.pid \
		--exec $DAEMON
		echo done
	else
		echo disabled 
	fi
	;;
  stop)
	echo -n "Stopping $DESC: "

	start-stop-daemon --stop --oknodo --quiet --pidfile @localstatedir@/run/$NAME.pid \
		--exec $DAEMON
	echo done 
	$FIREWALL stop
	;;
  restart|force-reload)
	echo -n "Restarting $DESC: "

	start-stop-daemon --stop --oknodo --quiet --pidfile \
		@localstatedir@/run/$NAME.pid --exec $DAEMON
	sleep 1

	$FIREWALL restart

	start-stop-daemon --start --quiet --pidfile \
		@localstatedir@/run/$NAME.pid --exec $DAEMON
	echo "$NAME."
	;;
  *)
	N=@sysconfdir@/init.d/$NAME
	echo "Usage: $N {start|stop|restart|force-reload}" >&2
	exit 1
	;;
esac

exit 0
